# packaging/openrc/init.d/coolercontrol
#!/sbin/openrc-run

description="CoolerControl Daemon"
supervisor="supervise-daemon"
command="/usr/bin/coolercontrold"
command_args=""
COOLERCONTROL_LOG_LEVEL="INFO"
supervise_daemon_args=" -e "COOLERCONTROL_LOG=${COOLERCONTROL_LOG_LEVEL}""
output_logger="logger -et '${RC_SVCNAME}'"
error_logger="logger -et '${RC_SVCNAME}' -p3"

depend() {
    after net coolercontrol-liqctld 
    want coolercontrol-liqctld
    use logger
}
